# deploy/systemd/napcat-login-watchdog-venv.service
[Unit]
Description=NapCat login watchdog
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/napcat-login-watchdog
EnvironmentFile=/opt/napcat-login-watchdog/.env
ExecStart=/opt/napcat-login-watchdog/.venv/bin/napcat-login-watchdog run

